# timers/monerod-build.timer
[Unit]
Description=monerod build timer

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target